# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:squeezeboxserver]
command=squeezeboxserver --user root --prefsdir /mnt/state/prefs --logdir /mnt/state/logs --cachedir /mnt/state/cache

